# deploy/fx-model-api.service
[Unit]
Description=FX Model API (TFT live prediction service)
After=network-online.target mongod.service
Wants=network-online.target
Requires=mongod.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/fx-model-api/
Environment=PATH=/opt/fx-model-api/.env
# Uncomment to pin a GPU:
# Environment=CUDA_VISIBLE_DEVICES=0
ExecStart=/opt/fx-model-api/.env/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
# Give the model time to load before considering a restart a failure:
TimeoutStartSec=120
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
